# energy_dashboard/config/supervisor/energy_dashboard.conf
[program:energy_dashboard]
command=/home/energy_dashboard/.virtualenvs/energy_dashboard/bin/gunicorn energy_dashboard.wsgi:application --workers 2 --bind unix:/tmp/energy_dashboard.sock
user=energy_dashboard
directory=/home/energy_dashboard/energy_dashboard/
stdout_logfile=/var/log/energy_dashboard/gunicorn
redirect_stderr=true
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8
