# supervisor/energy_dashboard.conf
[program:energy_dashboard]
command=/home/energy_dashboard/.virtualenvs/energy_dashboard/bin/gunicorn energy_dashboard.wsgi:application
directory=/home/energy_dashboard/energy_dashboard
user=energy_dashboard
autostart=true
autorestart=true
redirect_stderr=true
